# logstash/cookbooks/is_logstash/templates/default/logstash.service.erb
[Unit]
Description=logstash

[Service]
Type=simple
User=<%= node['user'] %>
Group=<%= node['user'] %>
# Load env vars from /etc/default/ and /etc/sysconfig/ if they exist.
# Prefixing the path with '-' makes it try to load, but if the file doesn't
# exist, it continues onward.
EnvironmentFile=-/etc/default/logstash
EnvironmentFile=-/etc/sysconfig/logstash
ExecStart=/logstash-<%= node['version'] %>/bin/logstash "--path.settings" "/etc/logstash"
Restart=always
WorkingDirectory=/
Nice=19
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
